# xtree.service
[Unit]
Description=display.py

After=default.target

[Service]
ExecStart=/usr/local/bin/display.py

[Install]
WantedBy=default.target
